# scripts/worker-update.service
[Unit]
Description=WorkerUpdate

[Service]
Type=simple
User=vagrant
ExecStart=/vagrant/scripts/worker_update.sh